# files/etc/uci-defaults/10-custom-settings
#!/bin/sh
# Custom settings for Hostname, Description, Notes, and Wi-Fi configuration

# 1. Set Hostname, Description, and Notes
uci set system.@system[0].hostname='R1S-H5'
uci set system.@system[0].description='NanoPi R1S H5 Custom Router (by CrueChan)'
uci set system.@system[0].notes='Compiled on BUILD_DATE by CrueChan (ImmortalWrt 23.05.7, Kernel 5.15, Stable Release)'
uci commit system
# Apply hostname change immediately to current session
echo "R1S-H5" > /proc/sys/kernel/hostname

# 2. Set Default LAN IP to 192.168.2.1
uci set network.lan.ipaddr='192.168.2.1'
uci commit network

# 3. Set Default Root Password to 'password'
if [ -x /usr/sbin/chpasswd ]; then
    echo "root:password" | chpasswd
elif [ -x /usr/bin/chpasswd ]; then
    echo "root:password" | chpasswd
else
    # Fallback to interactive input redirection for passwd command
    (echo "password"; sleep 1; echo "password") | passwd root
fi

# 4. Configure UPnP (Enabled by default)
if [ -f /etc/config/upnpd ]; then
    uci set upnpd.config.enabled='1'
    uci set upnpd.config.enable_upnp='1'
    uci set upnpd.config.enable_natpmp='1'
    uci commit upnpd
fi

# 5. Configure Wireless defaults depending on version
# First, ensure /etc/config directory exists
mkdir -p /etc/config

# If the template wireless config is missing (e.g. after a factory reset on ext4), restore it.
if [ ! -f /etc/config/wireless ]; then
    cat << 'EOF' > /etc/config/wireless
config wifi-device 'radio0'
	option type 'mac80211'
	option hwmode '11g'
	option channel '11'
	option htmode 'HT20'
	option country 'AU'
	option txpower '15'
	option disabled '1'

config wifi-iface 'default_radio0'
	option device 'radio0'
	option network 'lan'
	option mode 'ap'
	option ssid 'FriendlyWrt'
	option encryption 'psk2'
	option key 'password'
	option disabled '1'
EOF
fi

# Load the wireless driver explicitly and wait for the wlan0 interface to register
modprobe 8189es 2>/dev/null
modprobe rtl8189es 2>/dev/null
for i in 1 2 3 4 5; do
    [ -d /sys/class/net/wlan0 ] && break
    sleep 1
done

# Resolve the exact sysfs path suffix of the wireless hardware dynamically
DEV_PATH=$(readlink -f /sys/class/net/wlan0/device 2>/dev/null)
PATH_SUFFIX=$(echo "$DEV_PATH" | sed 's|^/sys/devices/||')

if [ -n "$PATH_SUFFIX" ]; then
    uci set wireless.radio0.path="$PATH_SUFFIX"
fi

SYS_VERSION=$(cat /etc/image_version 2>/dev/null)

if [ "$SYS_VERSION" = "21.02.7" ]; then
    # For legacy 21.02: Enable both wireless device and AP interface
    uci set wireless.radio0.disabled='0'
    uci set wireless.default_radio0.disabled='0'
    
    # Configure AP SSID dynamically based on MAC address
    raw_mac=$(cat /sys/class/net/wlan0/address 2>/dev/null | tr -d ':' | tr 'a-z' 'A-Z')
    wlan_mac=$(echo "$raw_mac" | tr -d ' \t\r\n')
    if [ -n "$wlan_mac" ] && [ ${#wlan_mac} -ge 4 ]; then
        mac_suffix=$(echo -n "$wlan_mac" | tail -c 4)
    else
        mac_suffix="21.02"
    fi
    uci set wireless.default_radio0.ssid="FriendlyWrt-$mac_suffix"
    uci set wireless.default_radio0.encryption='psk2'
    uci set wireless.default_radio0.key='password'
else
    # For 23.05: Enable radio0 for scan/client, and disable the AP interface to prevent panic
    uci set wireless.radio0.disabled='0'
    uci set wireless.default_radio0.disabled='1'
fi

uci commit wireless

exit 0
