# files/etc/uci-defaults/10-custom-settings
#!/bin/sh
# Custom settings for Hostname, Description, Notes, and Wi-Fi configuration

# 1. Set Hostname, Description, and Notes
uci set system.@system[0].hostname='R1S-H5'
uci set system.@system[0].description='NanoPi R1S H5 Custom Router (by CrueChan)'
uci set system.@system[0].notes='Compiled on BUILD_DATE by CrueChan (ImmortalWrt 23.05.7, Kernel 5.15, Stable Release)'
uci commit system
# Apply hostname change immediately to current session
echo "R1S-H5" > /proc/sys/kernel/hostname

# 2. Set Default LAN IP to 192.168.2.1
uci set network.lan.ipaddr='192.168.2.1'
uci commit network

# 3. Set Default Root Password to 'password'
if [ -x /usr/sbin/chpasswd ]; then
    echo "root:password" | chpasswd
elif [ -x /usr/bin/chpasswd ]; then
    echo "root:password" | chpasswd
else
    # Fallback to interactive input redirection for passwd command
    (echo "password"; sleep 1; echo "password") | passwd root
fi

# 4. Configure UPnP (Enabled by default)
if [ -f /etc/config/upnpd ]; then
    uci set upnpd.config.enabled='1'
    uci set upnpd.config.enable_upnp='1'
    uci set upnpd.config.enable_natpmp='1'
    uci commit upnpd
fi

# 5. Configure Wireless defaults depending on version
# On 21.02 we regenerate the config using wifi detect to get the correct Kernel 5.4 device paths.
# On 23.05 we use the pre-created config file to skip wifi detect and prevent early boot probing crash.
mkdir -p /etc/config
SYS_VERSION=$(cat /etc/image_version 2>/dev/null)

if [ "$SYS_VERSION" = "21.02.7" ]; then
    # Load the driver explicitly to make sure it is registered before detection
    modprobe 8189es 2>/dev/null
    modprobe rtl8189es 2>/dev/null
    sleep 1
    
    # Regenerate config for 21.02 to bind to kernel 5.4 paths
    rm -f /etc/config/wireless
    /sbin/wifi detect > /etc/config/wireless
    
    if [ -s /etc/config/wireless ]; then
        uci set wireless.radio0.disabled='0'
        
        # Configure default AP if it was created, or add one if it doesn't exist
        if ! uci get wireless.@wifi-iface[0] >/dev/null 2>&1; then
            uci add wireless wifi-iface
            uci rename wireless.@wifi-iface[-1]='default_radio0'
            uci set wireless.default_radio0.device='radio0'
            uci set wireless.default_radio0.network='lan'
            uci set wireless.default_radio0.mode='ap'
        fi
        
        uci set wireless.@wifi-iface[0].disabled='0'
        
        # Configure AP SSID dynamically
        raw_mac=$(cat /sys/class/net/wlan0/address 2>/dev/null | tr -d ':' | tr 'a-z' 'A-Z')
        wlan_mac=$(echo "$raw_mac" | tr -d ' \t\r\n')
        if [ -n "$wlan_mac" ] && [ ${#wlan_mac} -ge 4 ]; then
            mac_suffix=$(echo -n "$wlan_mac" | tail -c 4)
        else
            mac_suffix="21.02"
        fi
        uci set wireless.@wifi-iface[0].ssid="FriendlyWrt-$mac_suffix"
        uci set wireless.@wifi-iface[0].encryption='psk2'
        uci set wireless.@wifi-iface[0].key='password'
        uci commit wireless
    fi
else
    # For 23.05: Keep pre-created configuration to avoid crash. 
    # If the file is missing (e.g. after a factory reset), restore the safe template first.
    if [ ! -f /etc/config/wireless ]; then
        cat << 'EOF' > /etc/config/wireless
config wifi-device 'radio0'
	option type 'mac80211'
	option path 'platform/soc/1c10000.mmc/mmc_host/mmc1/mmc1:0001/mmc1:0001:1'
	option hwmode '11g'
	option channel '11'
	option htmode 'HT20'
	option country 'AU'
	option txpower '15'
	option disabled '0'

config wifi-iface 'default_radio0'
	option device 'radio0'
	option network 'lan'
	option mode 'ap'
	option ssid 'FriendlyWrt'
	option encryption 'psk2'
	option key 'password'
	option disabled '1'
EOF
    fi
    # Enable radio0 for scan/client, and disable the AP interface
    uci set wireless.radio0.disabled='0'
    uci set wireless.default_radio0.disabled='1'
    uci commit wireless
fi

exit 0
