# files/etc/uci-defaults/10-custom-settings
#!/bin/sh
# Custom settings for Hostname, Description, Notes, and Wi-Fi configuration

# 1. Set Hostname, Description, and Notes
uci set system.@system[0].hostname='R1S-H5'
uci set system.@system[0].description='NanoPi R1S H5 Custom Router (by CrueChan)'
uci set system.@system[0].notes='Compiled on BUILD_DATE by CrueChan (ImmortalWrt 23.05.7, Kernel 5.15, Stable Release)'
uci commit system
# Apply hostname change immediately to current session
echo "R1S-H5" > /proc/sys/kernel/hostname

# 2. Safely disable onboard Wi-Fi by default to prevent kernel panic / boot loop
if [ -f /etc/config/wireless ]; then
    # Explicitly disable radio0 (RTL8189ES) on boot
    uci set wireless.radio0.disabled='1'
    uci commit wireless
    
    # Apply and stop any active wireless interface
    wifi reload
fi

exit 0
